# host-agent/hostagent.service
[Unit]
Description=S3 Custom Host Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env node /opt/s3custom/host-agent/dist/index.js
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/home/jirka/.local/share/mise/installs/node/25.2.1/bin
Restart=on-failure
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
